# systemd/agent.service
[Unit]
Description=Multi-agent AI daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/agent
ExecStart=/home/ubuntu/.local/bin/uv run python -m agent.main --no-cli
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/ubuntu/.nvm/versions/node/v24.14.0/bin:/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target